# root/etc/services.d/repertory/run
#!/usr/bin/with-contenv bash

cd /app/repertory || exit

echo "Starting Repertory."

exec \
	s6-setuidgid abc /bin/bash -c
